# common/systemd/.config/systemd/user/weather-update.service
[Unit]
Description=Update the weather cache file (JSON) fetched from wttr.in

[Service]
Type=oneshot
ExecStart=%h/.config/wttrin/weather-update
